# services/social/systemd/mochirii-social-backup.timer
[Unit]
Description=Run Mochirii Social backup nightly at 03:15 UTC

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
AccuracySec=1min
Persistent=true
Unit=mochirii-social-backup.service

[Install]
WantedBy=timers.target
